# scripts/xyz.boarlabs.slack-heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>xyz.boarlabs.slack-heartbeat</string>

    <key>WorkingDirectory</key>
    <string>/Users/peterovchinnikov/cron_workers/linkedin-ai/worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/peterovchinnikov/cron_workers/linkedin-ai/worker/scripts/cron-wrapper.sh</string>
        <string>/Users/peterovchinnikov/cron_workers/linkedin-ai/worker/scripts/slack-heartbeat.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/peterovchinnikov/cron_workers/linkedin-ai/logs/slack-heartbeat.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/peterovchinnikov/cron_workers/linkedin-ai/logs/slack-heartbeat.err.log</string>

    <key>SessionCreate</key>
    <false/>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
