# scripts/xyz.boarlabs.slack-heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>xyz.boarlabs.slack-heartbeat</string>

    <key>WorkingDirectory</key>
    <string>/Users/peterovchinnikov/cron_workers/linkedin-ai/worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/peterovchinnikov/cron_workers/linkedin-ai/worker/scripts/cron-wrapper.sh</string>
        <string>/Users/peterovchinnikov/cron_workers/linkedin-ai/worker/scripts/slack-heartbeat.sh</string>
    </array>

    <!--
      Schedule (system-local time; this Mac is Europe/Podgorica): once a day
      at 06:00, Tuesday through Friday only. Each entry carries a Weekday so
      launchd scopes it to that day (0/7=Sun, 1=Mon, 2=Tue, 3=Wed, 4=Thu,
      5=Fri, 6=Sat):
        Sat, Sun, Mon   - no fires
        Tue-Fri (2-5)   - 06:00
      If the Mac is asleep at 06:00, launchd coalesces missed events into one
      fire on wake; a powered-off Mac gets no catch-up fire.
    -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <!-- Grace between launchd's SIGTERM and SIGKILL when a running fire is
         force-stopped (bootout / kickstart -k / logout). 60 is the launchd
         LaunchAgent ceiling — higher values are silently clamped to 60
         (verified 2026-07-17: 45→45, 150→60) — so it's the max window the
         run-finished Slack bookend gets to post on a forced stop. Enough for
         the typical sub-30s haiku call; normal completions don't hit this. -->
    <key>ExitTimeOut</key>
    <integer>60</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/peterovchinnikov/cron_workers/linkedin-ai/logs/slack-heartbeat.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/peterovchinnikov/cron_workers/linkedin-ai/logs/slack-heartbeat.err.log</string>

    <key>SessionCreate</key>
    <false/>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
